# deploy/supervisord.ini
[program:fics_web_interface]
directory=/srv/fics
command=/usr/local/bin/node lib/server.js
numprocs=1
process_name=fics_web_interface
autostart=true
autorestart=true
user=apache
environment=PATH="/usr/local/bin/"
