# deploy/supervisord.ini
[program:fics_web_interface]
directory=/srv/fics
command=/usr/local/bin/npm start
numprocs=1
process_name=%(program_name)
autostart=true
autorestart=true
user=apache
environment=PATH="/usr/local/bin/"
